# meta-rdk-camera/recipes-multimedia/cvr/files/cvr-launcher.service
[Unit]
Description=CVR
After=multi-user.target pipewire-launcher.service
Requires=pipewire-launcher.service

[Service]
ExecStart=/bin/sh -c '(/lib/rdk/startCVR.sh)' &
Type=forking
Restart=always

[Install]
WantedBy=multi-user.target
